# client/manager/static/Artisan_Manager.service
[Unit]
Description=Artisan Hosting Machine Management System
After=network.target
Wants=ArtisanManager.timer

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/env python3 /usr/local/bin/Artisan/bin/manage_machine

[Install]
WantedBy=timers.target
